# vagrant/upstart-private-ledger.conf
start on vagrant-mounted
stop on runlevel [06]

respawn

exec /stellard-src/build/stellard --conf /stellard-src/vagrant/stellard-private-ledger.cfg --load
