# vagrant/upstart-public-ledger.conf
start on vagrant-mounted
stop on runlevel [06]
stop on stellard-reprovision

respawn

exec /stellard-src/build/stellard --conf /stellard-src/vagrant/stellard-public-ledger.cfg
