# vagrant/upstart.conf
start on runlevel [2345]
stop on runlevel [06]

respawn

exec /stellard-src/build/stellard --conf /stellard-src/vagrant/stellard.cfg --load
